# launchd/com.agentic.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/g4ndr1k</string>
		<key>MallocStackLogging</key>
		<string>0</string>
		<key>PYTHONPATH</key>
		<string>/Users/g4ndr1k/agentic-ai</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.agentic.bridge</string>
	<key>SoftResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>8192</integer>
	</dict>
	<key>HardResourceLimits</key>
	<dict>
		<key>NumberOfFiles</key>
		<integer>65536</integer>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/AgenticAI.app/Contents/MacOS/launch_bridge</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/g4ndr1k/agentic-ai/logs/bridge-launchd-err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/g4ndr1k/agentic-ai/logs/bridge-launchd.log</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>WorkingDirectory</key>
	<string>/Users/g4ndr1k/agentic-ai</string>
</dict>
</plist>
